# Docs/orleans-edge-node.service
[Unit]
Description=Orleans Edge Node - IoT LED Controller Silo
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/orleans-edge
ExecStart=/home/pi/.dotnet/dotnet OrleansEdge.Node.dll
Restart=always
RestartSec=10
TimeoutStartSec=300
SyslogIdentifier=orleans-edge-node

# Environment variables
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production

# Security hardening (optional but recommended)
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
